# custom_services/service-weather.service
[Unit]
Description=Weather Service

[Service]
ExecStart=/usr/local/bin/weather.sh
Type=oneshot

[Install]
WantedBy=default.target
